# services/sphinxsearch/crontab/root
*/5 * * * * indexer --rotate --all >> /dev/null 2>&1
